# fhgfs_admon/build/dist/usr/lib/systemd/system/beegfs-admon.service
[Unit]
Description=Start BeeGFS Admon
Requires=network-online.target
After=network-online.target zfs.target

[Service]
ExecStart=/etc/init.d/beegfs-admon start
ExecStop=/etc/init.d/beegfs-admon stop
Type=forking
# We disable PIDFile= because it doesn't work with multi-mode configurations
#PIDFile=/var/run/beegfs-admon.pid

[Install]
WantedBy=multi-user.target