# event_listener/build/dist/usr/lib/systemd/system/beegfs-eventlistener.service
[Unit]
Description=BeeGFS event-listener Server
Documentation=https://doc.beegfs.io/latest/
After=network-online.target

[Service]
Environment=CONFIG_FILE=/etc/beegfs/beegfs-eventlistener.conf
ExecStart=/etc/init.d/beegfs-eventlistener $CONFIG_FILE
Type=simple

StandardOutput=append:/var/log/beegfs-eventlistener.log
StandardError=append:/var/log/beegfs-eventlistener.log

[Install]
WantedBy=multi-user.target
