# helperd/build/dist/usr/lib/systemd/system/beegfs-helperd.service
[Unit]
Description=BeeGFS Helperd
Documentation=http://www.beegfs.com/content/documentation/
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/beegfs/sbin/beegfs-helperd cfgFile=/etc/beegfs/beegfs-helperd.conf runDaemonized=false
Type=simple

[Install]
WantedBy=multi-user.target
